# puppet-example-ce/modules/multi_server/templates/site_instance/codenvy.service.erb
[Unit]
Description=Codenvy Site Server
After=network.target

[Service]
Type=forking
ExecStart=/usr/lib/systemd/codenvy start
ExecStop=/usr/lib/systemd/codenvy stop
PIDFile=/home/codenvy/codenvy.pid

[Install]
WantedBy=multi-user.target
