# puppet-example-ce/modules/third_party/templates/gorouter/gorouter.service.erb
[Unit]
Description=gorouter
After=network.target

[Service]
Type=forking
ExecStart=/usr/lib/systemd/gorouter start
ExecStop=/usr/lib/systemd/gorouter stop
PIDFile=/var/run/gorouter.pid

[Install]
WantedBy=multi-user.target
